# nyt-express-backend/scripts/crontab
# Run puzzle refresh daily at 5 AM UTC
0 5 * * * cd /path/to/nyt-crossword-leaderboard/nyt-express-backend && ./scripts/refresh-puzzle.sh 